# server/babiling.socket
[Unit]
Description=Activation socket for Babiling

[Socket]
ListenStream=3468

[Install]
WantedBy=sockets.target
